# packages/nodewatcher-agent/files/nodewatcher-agent.init
#!/bin/sh /etc/rc.common

START=90

USE_PROCD=1
NAME=nodewatcher-agent
PROG=/usr/bin/nodewatcher-agent

start_service() {
	procd_open_instance
	procd_set_param command "$PROG"
	procd_close_instance
}

stop() {
	service_stop "$PROG"
}

reload() {
	service_reload "$PROG"
}
