# packages/serval-dna/files/etc/uci-defaults/serval-dna
#!/bin/sh

echo "export SERVALINSTANCE_PATH=/etc/serval" >> /etc/profile

/usr/bin/servald keyring add

exit 0
