# packages/serval-dna/files/etc/init.d/serval-dna
#!/bin/sh /etc/rc.common

START=93

boot() {
  /usr/bin/servald start
}

start() {
  boot
}
