# packages/serval-dna/files/etc/init.d/serval-dna
#!/bin/sh /etc/rc.common

START=63

boot() {
  [ -z `SERVALINSTANCE_PATH=/etc/serval /usr/bin/servald keyring list` ] && SERVALINSTANCE_PATH=/etc/serval /usr/bin/servald keyring add
  SERVALINSTANCE_PATH=/etc/serval /usr/bin/servald start
}

start() {
  boot
}
